# podman/openmemory.container
[Unit]
Description=OpenMemory - Hierarchical Semantic Graph Memory System
Documentation=https://github.com/nullure/openmemory
After=network-online.target
Wants=network-online.target

[Container]
Image=ghcr.io/nullure/openmemory:latest
ContainerName=openmemory
Network=pasta
PublishPort=8080:8080
Volume=openmemory-data.volume:/data:Z
EnvironmentFile=%h/.config/openmemory/openmemory.env
User=1001:1001
HealthCmd=bun -e "fetch('http://localhost:8080/health').then(r=>process.exit(r.ok?0:1)).catch(()=>process.exit(1))"
HealthInterval=30s
HealthTimeout=10s
HealthRetries=3
HealthStartPeriod=30s
Restart=always
RestartSec=10s

[Service]
TimeoutStartSec=900
TimeoutStopSec=30
Restart=on-failure
RestartSec=10s

[Install]
WantedBy=default.target
